# usb_hid_gadget.service
[Unit]
Description=USB HID Gadget Service

[Service]
Type=oneshot
ExecStart=/home/dottodots/ai-braille-keyboard/usb-hid-gadget.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
